# .claude/infra/tovplay-watchdog.service
[Unit]
Description=TovPlay Docker Watchdog Service
Documentation=https://github.com/tovtech/tovplay
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Environment
Environment=SLACK_WEBHOOK_URL=
Environment=CHECK_INTERVAL=30

# Execution
ExecStart=/opt/tovplay/scripts/docker-watchdog.sh
ExecReload=/bin/kill -HUP $MAINPID

# Restart behavior - always restart no matter what
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Resource limits
MemoryMax=128M
CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tovplay-watchdog

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
